# system/roles/service-containers/files/chat/templates/ircd/supervisord.conf.j2
[supervisord]
nodaemon=true

[program:ircd]
command=/usr/sbin/ircd-hybrid -foreground
user=irc
